# sv/packetbeat/run
#!/bin/sh

sv start elasticsearch || exit 1

cd /packetbeat

exec 2>&1
exec ./packetbeat -c ./packetbeat.yml
