# sv/elasticsearch/run
#!/bin/sh

sed -i -e "s|#cluster.name:.*|cluster.name: `hostname`|" /elasticsearch/config/elasticsearch.yml

exec 2>&1
exec /elasticsearch/bin/elasticsearch
