# sv/kibana/run
#!/bin/bash

sv start elasticsearch || exit 1

source /etc/envvars
env

curl -XGET "http://localhost:9200/_cluster/health?wait_for_status=green&timeout=30s"

exec 2>&1
exec /kibana/bin/kibana
