# sv/kibana/run
#!/bin/sh

sv start elasticsearch || exit 1
curl -XGET 'http://localhost:9200/_cluster/health?wait_for_status=green&timeout=30s'

cd /kibana/bin

exec 2>&1
exec ./kibana
